# conf/systemd/reservoirdog-prod.service
[Unit]
Description=Reservoir Dog - Production Environment
After=network.target

[Service]
Type=simple
User=ansible
Group=ansible
WorkingDirectory=/opt/prod/reservoirdog
Environment="ENVIRONMENT=prod"
Environment="FLASK_PORT=45080"
ExecStart=/opt/prod/reservoirdog/venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

